# api_openweathermap/cron/api_openweathermap.crontab
EMAILS="huangpj1110@gmail.com"
PYTHONPATH=/anaconda2/lib/python2.7/site-packages

15 0-23/8 * * * python ~/dev/api/api_openweathermap/weather_current.py > ~/log/weather_current.log 2>&1; cat ~/log/weather_current.log | mail -s "Weather Report cron output" "$EMAILS"